# cookbooks/baasbox/templates/default/baasbox.conf.erb
start on runlevel [2345]
stop on runlevel [06]

script
  cd /opt/baasbox;
  exec ./start > /var/log/baasbox.log
end script
